# 亚马逊补货/deploy/amazon-restock.service
[Unit]
Description=Amazon Restock Analysis Service
Documentation=https://github.com/your-repo/amazon-restock
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/path/to/your/project
Environment=PATH=/path/to/your/project/.venv/bin
Environment=PYTHONPATH=/path/to/your/project
ExecStart=/path/to/your/project/.venv/bin/python /path/to/your/project/main.py --server
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amazon-restock

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/your/project/logs
ReadWritePaths=/path/to/your/project/output
ReadWritePaths=/path/to/your/project/data

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 